# etc/psmsd.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<!--
	Copyright (c) 2016 Peter Eriksson (peter@lysator.liu.se)

	This manifest can be used to manage the PSMSD daemon using the 
	Solaris SMF subsystem.

	Import this manifest using: 

		svccfg -v import psmsd.xml
	
        Then activate the daemon with:

		svcadm enable /site/psmsd:default
-->

<service_bundle type='manifest' name='psmsd'>

<service
	name='site/psmsd'
	type='service'
	version='1'>

	<dependency
		name='fs-local'
		type='service'
		grouping='require_all'
		restart_on='none'>
		<service_fmri value='svc:/system/filesystem/local' />
	</dependency>

        <exec_method
                type='method'
                name='stop'
                exec=':kill -9'
                timeout_seconds='60'>
        </exec_method>

        <instance name='default' enabled='false'>

        <exec_method
                type='method'
                name='start'
                exec='/usr/local/sbin/start-psmsd'
                timeout_seconds='60'>
        </exec_method>

	</instance>

	<stability value='Evolving' />

	<template>
		<common_name>
			<loctext xml:lang='C'>
			PSMSD Server
			</loctext>
		</common_name>
		<documentation>
			<manpage title='psmsd' section='1' manpath='/usr/local/share/man' />
		</documentation>
	</template>
</service>

</service_bundle>
